# slave/files/mnt-jenkins.automount
[Unit]
Before=remote-fs.target

[Automount]
Where=/mnt/jenkins

[Install]
WantedBy=multi-user.target
