# slave/files/mnt-jenkins.mount
[Mount]
What=otc-mesa-ci.local:/srv/jenkins
Where=/mnt/jenkins
Type=nfs
Options=defaults
